# bin/rc.cmk-api.EXAMPLE
#!/bin/bash
#
# Init script for the cmk-api service
#
### BEGIN INIT INFO
# Provides: 
# Required-Start: $network
# Required-Stop: 
# Should-Start: 
# Should-Stop: 
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: start and stop cmk-api
# Description: cmk-api is a REST API for check_mk
### END INIT INFO

SITE="watotest"
PREFIX="/omd/sites/$SITE/cmk-api"
PIDFILE="$PREFIX/rack.pid"
PORT=5005

if [ `whoami` = $SITE ] ; then
  SUDO=""
else
  SUDO="sudo -u $SITE"
fi

case $1 in
start)
	printf "Starting cmk-api for site $SITE.. "
        $SUDO scl enable ruby193 "bundle exec '$PREFIX/bin/rackup -p $PORT -D --pid $PREFIX/rack.pid'"
 	echo "done"
	;;

start-debug)
	echo "Starting cmk-api for site $SITE in debugging mode.. "
        $SUDO scl enable ruby193 "bundle exec '$PREFIX/bin/rackup -p $PORT'"
	;;

  
stop)
	pid=`cat $PIDFILE`
	printf "Stopping cmk-api for site $SITE (pid $pid).. "
	$SUDO kill -INT $pid
	rm $PIDFILE
	echo "done"
	;;

restart)
	echo "FIXME -- STUB"
	;;

*)
	echo "unknown action" 
	exit 1
esac

