# templates/rc.erb
#!/bin/bash
#
# Init script for the cmk-api service
#
### BEGIN INIT INFO
# Provides: 
# Required-Start: $network
# Required-Stop: 
# Should-Start: 
# Should-Stop: 
# Default-Start: 3 4 5
# Default-Stop: 0 1 2 6
# Short-Description: start and stop cmk-api
# Description: cmk-api is a REST API for check_mk
### END INIT INFO

# Determine the site name by looking at the name of this script
# Example:   /etc/init.d/cmk-api.foo is for the 'foo' site
SITE=`echo $0 | sed 's/.*\.//'`

# Set the defaults
PREFIX="/opt/cmk-api"
PIDFILE="/omd/sites/$SITE/tmp/run/cmk-api.pid"
RUBY_BINDIR="/opt/ruby/2.0.0/bin"
RUN_AS="$SITE"
PORT=5005

# Read the site-specific overrides
SITE_CONFIG="/etc/sysconfig/cmk-api.$SITE"
if [ -f $SITE_CONFIG ] ; then
  . /etc/sysconfig/cmk-api.$SITE
fi

# Use sudo to switch to the site owner, if needed
if [ `whoami` != $RUN_AS ] ; then
  sudo -i -u $RUN_AS $0 $*
  exit $?
fi

# Switch to the preferred Ruby version
export PATH=$RUBY_BINDIR:$PATH

cd $PREFIX

case $1 in
start)
	printf "Starting cmk-api for site $SITE on port $PORT.. "
	bundle exec $PREFIX/bin/rackup -p $PORT -D --pid $PIDFILE
 	echo "done"
	;;

start-debug)
	echo "Starting cmk-api for site $SITE in debugging mode.. "
	bundle exec $PREFIX/bin/rackup -p $PORT
	;;

stop)
	pid=`cat $PIDFILE`
	printf "Stopping cmk-api for site $SITE (pid $pid).. "
	$SUDO kill -INT $pid
	rm $PIDFILE
	echo "done"
	;;

restart)
	$0 stop
	sleep 3
	$0 start
	;;

*)
	echo "unknown action" 
	exit 1
esac
